# openwrt/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

echo "nameserver 8.8.8.8" > /etc/resolv.conf
ping -c 4 8.8.8.8

/bin/sh /etc/config/watchdog.sh >> /dev/null 2>&1 &

while true; do
	/usr/bin/python /teplomer.py
	killall python
	sleep 1
	killall -9 python
	sleep 120
done


exit 0
